# systemd/snano-prune.service
[Unit]
Description=snano prune

[Service]
Type=oneshot
ExecStart=%h/bin/snano --prune
